# services/kanshi.service
[Unit]
Description=Kanshi - display output management
 
[Service]
Type=simple
ExecStart=/usr/bin/kanshi
 
[Install]
WantedBy=default.target
