# kmonad/kmonad@laptop-kb.path
[Unit]
Description=start kmonad for laptop keyboard
StartLimitIntervalSec=400
StartLimitBurst=3

[Install]
DefaultInstance=config
WantedBy=default.target
